# LinuxService/transcriptcleaner.service
[Unit]
Description=Transcript Cleaner
After=network.target
Documentation=https://thealmostengineer.com/

[Service]
Type=simple
Restart=always
ExecStart=/home/iamadmin/microservices/Almsotengr.TranscriptCleaner.Api
User=iamadmin

[Install]
WantedBy=multi-user.target